# deploy/workspace-backup.timer
[Unit]
Description=Daily workspace backup at 3am UTC

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
