# deploy/workspace-backup.service
[Unit]
Description=Back up ClawdBot workspace directories
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/backup-workspaces.sh
